# chef/site-cookbooks/selenium/templates/default/selenium-server.service.erb
[Unit]
Description=Selenium Server

[Service]
ExecStart=/bin/java -jar <%= File.join(node['selenium']['server']['install_path'], "selenium-server-standalone-#{node['selenium']['server']['version']}.jar") %>

[Install]
Alias=selenium-server.service
